# autorandr/90-display-smart-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/systemctl start --no-block display-smart-hotplug.service"
